# contrail/environment/modules/contrail/templates/contrail-api.svc.erb
#!/bin/sh

# chkconfig: 2345 99 01
# description: Juniper Network Virtualization API

<% for worker_id in (0..api_nworkers.to_i-1) %>
if [ -e /tmp/supervisord_config.sock ]; then
    supervisorctl -s unix:///tmp/supervisord_config.sock ${1} `basename ${0}:<%= worker_id %>`
else
    supervisorctl -s unix:///var/run/supervisord_config.sock ${1} `basename ${0}:<%= worker_id %>`
fi
<% end %>

