# contrail/environment/modules/contrail/templates/monitor_intf.service.erb
[Unit]
Description=Wait For INTF
Requisite=NetworkManager.service
After=NetworkManager.service
Wants=network.target rabbitmq-server.service
Before=network.target network-online.target

[Service]
Type=oneshot
ExecStart=/etc/monitor.sh

[Install]
WantedBy=multi-user.target

